# system/systemd/archivist/transmission.service.d/override.conf
[Unit]
After=
After=network.target Alexandria.mount
Requires=
Requires=Alexandria.mount
OnFailure=
OnFailure=status-email-user@%n.service status-send-telegram@%n.service
