# systemd/check_clean_boot.timer
[Unit]
Description=Check for unclean shutdown at boot

[Timer]
Unit=check_clean_boot.service
OnBootSec=0

[Install]
WantedBy=timers.target
